# etc/supervisor/conf.d/selenium.conf
[program:selenium-standalone-server]
command=java -jar /usr/local/bin/selenium-server-standalone.jar
autostart=true
autorestart=true