# etc/supervisor/conf.d/selenium.conf
[program:selenium-standalone-server]
command=java -jar selenium-standalone-server.jar
autostart=true
autorestart=true